# crontab
@reboot     ~/plant_light/plant_light.py auto
0  4 * * *  ~/plant_light/plant_light.py on
0 22 * * *  ~/plant_light/plant_light.py off
